# podman/systemd/mga-postgres.container
[Unit]
Description=MGA PostgreSQL 15 Database
Wants=network-online.target
After=network-online.target mga-network.network

[Container]
# PostgreSQL 15 official image (RHEL registry requires authentication)
Image=docker.io/library/postgres:15
ContainerName=mga-postgres
Network=mga-web
PublishPort=127.0.0.1:5432:5432
EnvironmentFile=%h/.config/mga-soap-calculator/postgres.env

# Persistent data with SELinux relabel (private mount)
# Official PostgreSQL image uses /var/lib/postgresql/data
Volume=mga-pgdata:/var/lib/postgresql/data:Z

# Health check - proper health verification
# Simplified command (RHEL image has native healthcheck support)
HealthCmd=pg_isready -U postgres -d mga_soap_calculator
HealthInterval=30s
HealthTimeout=10s
HealthRetries=5
HealthStartPeriod=30s

# Resource limits
Memory=1G

# Logging
PodmanArgs=--log-driver=journald

[Service]
Restart=always
RestartSec=5
TimeoutStartSec=300
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
