# podman/systemd/soap-calculator-api.container
[Unit]
Description=MGA Soap Calculator API Service
Wants=network-online.target mga-postgres.service
After=network-online.target mga-postgres.service mga-network.network
Requires=mga-postgres.service

[Container]
Image=localhost/mga-soap-calculator:latest
ContainerName=soap-api
Network=mga-web
PublishPort=8000:8000
EnvironmentFile=/home/grimm/.config/mga-soap-calculator/api.env

# Resource limits
Memory=1G

# Health check (readiness verification)
HealthCmd=/usr/bin/curl -fsS http://127.0.0.1:8000/api/v1/health || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=5
HealthStartPeriod=60s

# Logging integration
PodmanArgs=--log-driver=journald

[Service]
Restart=always
RestartSec=5
TimeoutStartSec=300
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
